# assets/postkeeper.service
[Unit]
Description=Postkeeper Milter

[Service]
ExecStart=/usr/sbin/postkeeper

[Install]
WantedBy=multi-user.target
